# inittab
tty1::respawn:/system/sh
tty2::respawn:/system/sh
tty3::respawn:/system/sh
tty4::respawn:/system/sh
# ::sysinit:printf "Welcome to Woven"
# ::restart:/sbin/init
# ::shutdown:printf -e "\nSyncing all file buffers."
# ::shutdown:sync
# ::shutdown:printf "Unmounting all filesystems."
# ::shutdown:umount -a -r
# ::ctrlaltdel:/sbin/reboot
